# lliurex-ssh-login.install/etc/update-motd.d/00-header.lliurex
#!/bin/sh

echo " "  
echo " __         __         __     __  __     ______     ______     __  __  "
echo "/\ \       /\ \       /\ \   /\ \/\ \   /\  == \   /\  ___\   /\_\_\_\   "
echo "\ \ \____  \ \ \____  \ \ \  \ \ \_\ \  \ \  __<   \ \  __\   \/_/\_\/_  "
echo " \ \_____\  \ \_____\  \ \_\  \ \_____\  \ \_\ \_\  \ \_____\   /\_\/\_\ " 
echo "  \/_____/   \/_____/   \/_/   \/_____/   \/_/ /_/   \/_____/   \/_/\/_/ "
echo ""

printf "Bienvenido a LliureX\n\n"
