# lliurex-ssh-login.install/etc/update-motd.d/10-help-text.lliurex
#!/bin/bash

URL="https://wiki.edu.gva.es/lliurex"

printf "\nDocumentación: %s\n\n" "$URL"
